# packages/sysutils/amremote/system.d/amremote.service
[Unit]
Description=Toggle between amremote and meson-ir
PartOf=lircd.service
After=lircd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/libreelec/meson_remote_ir
ExecReload=/usr/lib/libreelec/meson_remote_ir
RemainAfterExit=yes

[Install]
WantedBy=lircd.service
