# packages/sysutils/busybox/system.d/passwd-config.service
[Unit]
Description=Setup password config files
DefaultDependencies=no
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec /bin/sh /usr/lib/libreelec/passwd-config"
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
